# 332.conf
[program:chain33]
command=bash /www/chain33/chain33.sh
